# lib/desktop/customization-scripts/rc.local
#!/bin/sh -e
# remove the installer icon.
if [[ -f /home/lubuntu/Desktop/ubiquity.desktop ]]
then
  rm /home/lubuntu/Desktop/ubiquity.desktop
fi

exit 0
